# arris-s34-exporter.service
[Unit]
Description=Arris S34 scrape to node_exporter textfile collector
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/sysconfig/arris-s34-exporter
ExecStart=/usr/bin/arris_s34_write_prom.sh

[Install]
WantedBy=multi-user.target
